# client/openwrt/temp_files/crontab
*/10 * * * * /root/ssh_vpn_tun101/check-if-vpn-ok.sh
